# roles/containers/files/container2/root.crontab
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/4 * * * * > /var/log/apache2/access.log
* * * * * rm /var/www/html/uploads/*
